# macbook_scraper/com.vietmac.price-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vietmac.price-update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/gaurav/Documents/VietMac/macbook_scraper/auto-update-prices.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <!-- Run every 6 hours: 00:00, 06:00, 12:00, 18:00 -->
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/gaurav/Documents/VietMac/macbook_scraper/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/gaurav/Documents/VietMac/macbook_scraper/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/gaurav/Documents/VietMac/macbook_scraper</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
