# macbook_scraper/systemd/macbook-price-update.timer
[Unit]
Description=MacBook Price Update Timer
Requires=macbook-price-update.service

[Timer]
# Run daily at 2 AM
OnCalendar=daily
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
